# srcs/wordpress/supervisord.conf
[supervisord]
nodaemon=true

[program:php-fpm7]
command=/usr/sbin/php-fpm7 -D --fpm-config /etc/php7/php-fpm.d/www.conf  -R
autorestart=true

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autorestart=true

[eventlistener:processes]
command=sh -c "printf 'READY\n' && while read line; do kill -SIGQUIT $PPID; done < /dev/stdin"
events=PROCESS_STATE_STOPPED,PROCESS_STATE_EXITED,PROCESS_STATE_FATAL